# log/supervisor.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/bin/python /root/run_cron.py

[program:speedtest-web]
directory=/speedtest-easy/web
command=/usr/local/bin/node speedtest-web.js
